# deployment/systemd/alphapulse-backup.timer
[Unit]
Description=AlphaPulse-RL Backup Timer
Requires=alphapulse-backup.service

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target